# cathedral_orchestrator/rootfs/etc/services.d/cathedral/run
#!/command/with-contenv bash
exec /opt/app/run.sh
