# cathedral_orchestrator/rootfs/etc/services.d/cathedral/run
#!/command/execlineb -P
with-contenv
exec /opt/app/start.sh
